# deploy/systemd/zo-qore-fallback-watcher.service
[Unit]
Description=Zo-Qore Fallback Watcher
After=network.target zo-qore.service

[Service]
Type=simple
WorkingDirectory=/opt/zo-qore
Environment=NODE_ENV=production
EnvironmentFile=-/etc/zo-qore/env
ExecStart=/usr/bin/node /opt/zo-qore/dist/zo/fallback/start-watcher.js
Restart=always
RestartSec=5
User=zoqore
Group=zoqore

[Install]
WantedBy=multi-user.target
